# scripts/systemd/hololive-main-log-mirror@.service
[Unit]
Description=Hololive remote logs into repo logs (%i)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/work/hololive-bot
EnvironmentFile=-/etc/hololive/log-mirror.env
Environment=LOG_ROOT=/root/work/hololive-bot/logs
Environment=HOL_LOG_OSAKA_SSH_KEY=/root/work/hololive-bot/KR.key
Environment=HOL_LOG_SEOUL_SSH_KEY=/root/work/hololive-bot/KR.key
ExecStart=/root/work/hololive-bot/scripts/logs/remote-sync-main-logs.sh once %i
UMask=0027
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=6
