# scripts/systemd/hololive-compose.service
[Unit]
Description=hololive-bot Docker Compose stack
After=docker.service tailscaled.service network-online.target openbao-agent-hololive-bot.service
Requires=docker.service openbao-agent-hololive-bot.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/hololive-bot/compose/current
ExecStart=/usr/local/sbin/hololive-compose-up
ExecStop=/usr/local/sbin/hololive-compose-down
RemainAfterExit=yes
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
