# docs/hololive-bot-integrated-refactor-v3/scripts/systemd/hololive-main-log-mirror@.service
[Unit]
Description=Hololive remote logs into main /logs (%i)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=kapu
Group=docker
WorkingDirectory=/home/kapu/gemini/hololive-bot
EnvironmentFile=-/etc/hololive/log-mirror.env
Environment=LOG_ROOT=/logs
ExecStart=/home/kapu/gemini/hololive-bot/scripts/logs/remote-sync-main-logs.sh once %i
UMask=0027
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=6
