# scripts/systemd/hololive-compose.service.d/10-supplementary-groups.conf
[Service]
SupplementaryGroups=hololive-runtime docker
