# scripts/systemd/hololive-compose.service.d/20-live-compat.conf
[Service]
Environment=HOLOLIVE_ENABLE_LIVE_COMPAT=1
