# scripts/systemd/hololive-compose.service.d/30-ap-log-dir.conf
[Service]
Environment=HOL_AP_LOG_DIR=/home/kapu/work/iris-stack/hololive-bot/logs
